# conf/systemd/a2m.timer
[Unit]
Description=boot a2m

[Timer]
OnCalendar=Mon..Fri 14:00:00
Persistent=true

[Install]
WantedBy=timers.target
